# supervisord.conf
[program:node]
command = NODE_ENV=production node server.js
directory = /home/dotcloud/current