# supervisord.conf
[program:node]
environment = NODE_ENV=production